# Final-Evaluation/application/files/employee.service
[Unit]
Requires=network.target
After=network.target
[Service]
Type=simple
User=ubuntu
Group=ubuntu
Restart=always
RestartSec=3
Environment="CONFIG_FILE"="/opt/application/employee/config.yaml"
ExecStart=/opt/employee/employee
[Install]
WantedBy=multi-user.target
