# Final-Evaluation/application/files/frontend.service
[Unit]
Requires=network.target
After=network.target
[Service]
Type=simple
User=ubuntu
Group=ubuntu
Restart=always
RestartSec=3
WorkingDirectory=/opt/application/frontend
ExecStart=/usr/bin/serve -s build
[Install]
WantedBy=multi-user.target